# files/system/etc/profile.d/bun.sh

# These environment variables are the most reliable way to set global install paths for bun
export BUN_INSTALL_GLOBAL_DIR="/var/lib/bun/install/global"
export BUN_INSTALL_BIN="/usr/local/bin"
